# supervisor/conf.d/dlms.service.ini
[program:dlms.service]
command=/bin/bash -c 'source "$0" && exec "$@"' /etc/profile java -jar /opt/org.lowan.java/dlms.service/dlms.service-1.9.0.jar
user=root
autostart=true
autorestart=true
startsecs=1
startretries=999
redirect_stderr=false
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_events_enabled=true
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stderr_events_enabled=true
